# 50-load_news.rules
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idProduct}=="4227", ATTRS{idVendor}=="2237", RUN+="/usr/local/bin/send_signal.sh"
