# load_news.service
[Unit]
Description=load daily news on ereader
After=mnt-ereader.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/load_news.py
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
