# mnt-ereader.mount
[Unit]
Description=mount ereader

[Mount]
What=/dev/disk/by-uuid/xxx
Where=/mnt/ereader
Options=defaults
Type=vfat
TimeoutSec=10

[Install]
WantedBy=multi-user.target
